# random-wallpaper.timer
[Unit]
Description=random wallpaper timer

[Timer]
OnBootSec=0s
OnCalendar=*:0/1

[Install]
WantedBy=timers.target
