# random-wallpaper.service
[Unit]
Description=random wallpaper service

[Service]
Type=oneshot
EnvironmentFile=/etc/conf.d/random-wallpaper
ExecStart=/usr/bin/random-wallpaper
