# random-wallpaper.service
[Unit]
Description=random wallpaper

[Service]
Type=oneshot
ExecStart=/bin/bash /usr/local/bin/random-wallpaper.sh
