# random-wallpaper.service
[Unit]
Description=random wallpaper service

[Service]
Type=oneshot
EnvironmentFile=/etc/random-wallpaper.conf
ExecStart=/usr/bin/random-wallpaper
