# random-wallpaper.service
[Unit]
Description=random wallpaper service

[Service]
Type=oneshot
Environment=DISPLAY=:0
ExecStart=/usr/bin/random-wallpaper
